# ansible/roles/etcd-maintenance/templates/etcd-maintenance.service.j2
[Unit]
Description=Run etcd maintenance tasks
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/etcd-maintenance.sh
