# ansible/roles/etcd-maintenance/templates/etcd-maintenance.timer.j2
[Unit]
Description=Run etcd maintenance periodically

[Timer]
OnCalendar={{ etcd_maintenance_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target

